# root/etc/services.d/transmission/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh

export CURL_CA_BUNDLE=/etc/ssl/certs/ca-certificates.crt
s6-setuidgid twg_user transmission-daemon --foreground